# deploy/ilswlol.service.j2
[Unit]
Description=ilswlol

[Service]
Type=simple
ExecStart=/usr/bin/poetry run python runner.py
WorkingDirectory={{ site_path }}

[Install]
WantedBy=default.target
